# claude-code/rootfs/etc/s6-overlay/s6-rc.d/claude-code/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Claude Code Add-on: Service cleanup
# ==============================================================================

bashio::log.info "Claude Code terminal stopped (exit code: ${1})"
